# xiaozhi-server.service
[Unit]
Description=Xiaozhi ESP32 Server (Python)
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/xiaozhi-esp32-server/main/xiaozhi-server
ExecStart=/opt/conda/bin/conda run --no-capture-output -n xiaozhi-esp32-server python app.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xiaozhi-server

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
