# xiaozhi-manager-api.service
[Unit]
Description=Xiaozhi Manager API (Spring Boot)
After=network.target mysql.service redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/xiaozhi-esp32-server/main/manager-api
ExecStart=/usr/bin/java -Xms256m -Xmx512m -jar /opt/xiaozhi-esp32-server/main/manager-api/target/manager-api.jar --spring.profiles.active=dev
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xiaozhi-manager-api

[Install]
WantedBy=multi-user.target
